# pgAdmin4/etc/systemd/system/pgadmin.socket
[Unit]
Description=pgAdmin4 gunicorn socket

[Socket]
ListenStream=/run/pgadmin4.sock

[Install]
WantedBy=sockets.target
